# .config/autostart/unclutter.desktop
[Desktop Entry]
Name=Unclutter
Exec=unclutter --timeout 8
Terminal=false
Type=Application
X-GNOME-Autostart-enabled=true
